# gluster-lic-setup.in
#!/bin/bash
#
# chkconfig: 35 50 12
# description: Licensing kick-in
#


function main()
{
    if [ ! -d /lic ]; then
        mkdir -p /lic;
        cp @prefix@/ldpreload/gpgv /lic;
        cp @prefix@/ldpreload/pubring.gpg /lic;
    fi

    if [ ! -e /lib/libevil.so ]; then
        cp @prefix@/ldpreload/libevil.so /lib;
        chown root:root /lib/libevil.so;
        chmod g+s /lib/libevil.so;
    fi

    if [ ! -e /.epoch ]; then
        uuidgen -r > /.epoch;
    fi

    if [ ! -e /etc/glusterd/glusterd.info ]; then
        uuid=$(cat /.epoch);
        mkdir -p /etc/glusterd;
        echo "UUID=$uuid" > /etc/glusterd/glusterd.info;
    fi

    if [ ! -e /etc/ld.so.preload ]; then
        echo /lib/libevil.so > /etc/ld.so.preload;
    fi
}

main "$@";
